# package/DEBIAN/postinst
#!/bin/bash

systemctl daemon-reload

systemctl enable canary

systemctl start canary
